# template/home/config/autostart/autorun.desktop
[Desktop Entry]
Name=autorun.sh
Exec=/home/vagrant/.local/bin/autorun.sh
Terminal=false
Type=Application
