# template/home/config/autostart/create-rdp-shortcuts.desktop
[Desktop Entry]
Name=Create RDP connections shortcuts
Exec=/home/vagrant/.local/bin/create-rdp-shortcuts.sh
Terminal=false
Type=Application
