# template/home/config/autostart/endless-run-vmware-view.desktop
[Desktop Entry]
Name=endless-run-vmware-view.sh
Exec=/home/vagrant/.local/bin/endless-run-vmware-view.sh
Terminal=false
Type=Application
