# server/systemd/system/dnsmasq.service.d/override.conf
[Service]
Restart=always
RestartSec=1

[Unit]
# Disable start rate limiting
StartLimitInterval=0
