# template/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/mingetty --autologin {{ autologin_user }} --noclear %I
