# template/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=-/usr/bin/agetty --autologin ubuntu --noclear %I $TERM
