# template/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/mingetty --autologin ubuntu --noclear %I
